# .config/autostart/xfce4-panel.desktop
[Desktop Entry]
Type=Application
TryExec=xfce4-panel
Exec=xfce4-panel
Hidden=false
Name=Xfce Panel
